# deckrc/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root %I $TERM
ExecStop=/sbin/poweroff
Restart=no
StandardInput=tty
StandardOutput=tty
