# deckrc/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/startx
ExecStop=/sbin/poweroff
Restart=no
StandardInput=tty
StandardOutput=tty
